# crontab
#run python script every minutes
* * * * * /usr/local/bin/python /app/src/main.py
